# detect.service
[Service]
ExecStart = /home/sanveer/psl/psl_8/PSL8_TD/usbDetectServ.sh 
